# local/MACH/sup/crontab
30 2 * * 6,7	mckusick /usr/src/local/mach/mk/sup/machupd 2>&1 | mail mckusick
