# etc/SCCS/s.crontab
h43015
s 00000/00000/00013
d D 8.1 93/06/09 11:56:20 bostic 6 5
c 4.4BSD snapshot (revision 8.1)
e
s 00006/00000/00007
d D 5.5 93/06/09 11:56:16 bostic 5 4
c add header
e
s 00003/00003/00004
d D 5.4 90/06/30 21:58:06 bostic 4 3
c make into prototype
e
s 00003/00003/00004
d D 5.3 89/09/05 21:02:55 bostic 3 2
c mail daily/weekly/monthly output to root
e
s 00001/00001/00007
d D 5.2 89/04/17 22:14:01 bostic 2 1
c atrun -> /usr/libexec
e
s 00008/00000/00000
d D 5.1 89/04/12 12:17:12 bostic 1 0
c date and time created 89/04/12 12:17:12 by bostic
e
u
U
t
T
I 5
# Crontab file; see cron(8) for more information.
#
#	%W% (Berkeley) %G%

E 5
I 1
D 2
0,15,30,45 * * * *	root /usr/lib/atrun
E 2
I 2
0,15,30,45 * * * *	root /usr/libexec/atrun
I 5

E 5
E 2
D 4
00 1 * * *	daemon	/usr/new/lib/news/daily
0 0 * * *	uucp	/usr/local/lib/uucp.daily
7 5,12,18 * * *	uucp	/usr/local/lib/uucp.6hours
E 4
I 4
#0 1 * * *	daemon	/usr/contrib/news/daily
#0 0 * * *	uucp	/etc/uucp.daily
#7 5,12,18 * * *	uucp	/etc/uucp.6hours
I 5

E 5
E 4
D 3
0 2 * * *	root	/bin/sh /etc/daily 2>&1 | tee /var/log/daily.out | mail karels
30 3 * * 6	root	/bin/sh /etc/weekly >/var/log/weekly.out 2>&1
30 5 1 * *	root	/bin/sh /etc/monthly >/var/log/monthly.out 2>&1
E 3
I 3
0 2 * * *	root	/bin/sh /etc/daily 2>&1 | tee /var/log/daily.out | mail -s "daily output" root
30 3 * * 6	root	/bin/sh /etc/weekly 2>&1 | tee /var/log/weekly.out | mail -s "weekly output" root
30 5 1 * *	root	/bin/sh /etc/monthly 2>&1 | tee /var/log/monthly.out | mail -s "monthlyoutput" root
E 3
E 1
